# data/edu.ost.typetrace-backend.service.in
[Unit]
Description=TypeTrace Backend Service

[Service]
ExecStart=@bindir@/typetrace -b
Restart=on-failure

[Install]
WantedBy=default.target
